# service/motion/run
#!/bin/bash

sv start /service/motion-confd || exit 1
exec motion
